# conf/etc/services.d/01-postgres/run
#!/bin/ash
NAME="POSTGRES" 
COMMAND="/usr/local/bin/postgres-launcher" 
WAIT="/tmp/postgres"
exec service-handler "$NAME" "$COMMAND" "$WAIT"